# proj/djcelery.ini
[unix_http_server]
file=/home/dengsc/supervisor/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/home/dengsc/supervisor/supervisord.log  ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/home/dengsc/supervisor/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/dengsc/supervisor/supervisor.sock ; use a unix:// URL  for a unix socket

[program:djangoproject.celeryd]
command=/home/dengsc/gitapp/djcelery-example/venv/bin/python manage.py celeryd --concurrency=1 -E
user=dengsc
numprocs=1
directory=/home/dengsc/gitapp/djcelery-example
stdout_logfile=/home/dengsc/supervisor/log/celery_worker.log
stderr_logfile=/home/dengsc/supervisor/log/celery_worker.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 120
priority=998
environment=PYTHONPATH='$PYTHONPATH:/home/dengsc/gitapp/djcelery-example'
 
[program:djangoproject.celerybeat]
command=/home/dengsc/gitapp/djcelery-example/venv/bin/python manage.py celery beat --schedule=/tmp/celerybeat-schedule --pidfile=/tmp/django_celerybeat.pid --loglevel=INFO
user=dengsc
numprocs=1
directory=/home/dengsc/gitapp/djcelery-example
stdout_logfile=/home/dengsc/supervisor/log/celery_beat.log
stderr_logfile=/home/dengsc/supervisor/log/celery_beat.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 120
priority=998
environment=PYTHONPATH='$PYTHONPATH:/home/dengsc/gitapp/djcelery-example'
 
[program:djangoproject.celerycam]
command=/home/dengsc/gitapp/djcelery-example/venv/bin/python manage.py celerycam --frequency=10.0
user=dengsc
numprocs=1
directory=/home/dengsc/gitapp/djcelery-example
stdout_logfile=/home/dengsc/supervisor/log/celerycam.log
stderr_logfile=/home/dengsc/supervisor/log/celerycam.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 120
priority=998
environment=PYTHONPATH='$PYTHONPATH:/home/dengsc/gitapp/djcelery-example'
