# proj/djcelery.ini
[program:djangoproject.celeryd]
command=/home/dengsc/.local/share/virtualenvs/proj-aHI3p4-V/bin/python manage.py celeryd --concurrency=1 -E
user=dengsc
numprocs=1
directory=/home/dengsc/gitapp/proj
stdout_logfile=/var/log/celery_worker.log
stderr_logfile=/var/log/celery_worker.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 120
priority=998
environment=PYTHONPATH='$PYTHONPATH:/home/dengsc/gitapp/proj'
 
[program:djangoproject.celerybeat]
command=/home/dengsc/.local/share/virtualenvs/proj-aHI3p4-V/bin/python manage.py celery beat --schedule=/tmp/celerybeat-schedule --pidfile=/tmp/django_celerybeat.pid --loglevel=INFO
user=dengsc
numprocs=1
directory=/home/dengsc/gitapp/proj
stdout_logfile=/var/log/celery_beat.log
stderr_logfile=/var/log/celery_beat.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 120
priority=998
environment=PYTHONPATH='$PYTHONPATH:/home/dengsc/gitapp/proj'
 
[program:djangoproject.celerycam]
command=/home/dengsc/.local/share/virtualenvs/proj-aHI3p4-V/bin/python manage.py celerycam --frequency=10.0
user=dengsc
numprocs=1
directory=/home/dengsc/gitapp/proj
stdout_logfile=/var/log/celerycam.log
stderr_logfile=/var/log/celerycam.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 120
priority=998
environment=PYTHONPATH='$PYTHONPATH:/home/dengsc/gitapp/proj'
